# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:transactor]
command=/opt/datomic/bin/transactor /opt/datomic/config/transactor.properties
